# crux-ports/vdev/.footprint
drwxr-xr-x	t0b1nux/users	etc/
drwxr-xr-x	t0b1nux/users	etc/init.d/
-rwxr-xr-x	t0b1nux/users	etc/init.d/vdev
drwxr-xr-x	t0b1nux/users	etc/vdev/
drwxr-xr-x	t0b1nux/users	etc/vdev/actions/
-rw-r--r--	t0b1nux/users	etc/vdev/actions/001-sysfs-block.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/001-sysfs-char.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/002-modalias.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/003-firmware.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/agpgart.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/block.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/btrfs-control.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/char.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/cpu-msr.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/cpu.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/cuse.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/device-mapper.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/dialout.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/disk-part.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/disk.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/dm-disk.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/dri.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/err.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/fb.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/full.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/fuse-mounts.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/fuse.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/hwdb.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/ifname.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/input.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/irlpt.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/kmsg.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/kvm.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/legousbtower.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/lircN.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/loop-control.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/loop.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/lp.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/mISDNtimer.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/mem.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/mmtimer.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/mwave.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/net-tun.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/netdev.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/null.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/nvram.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/optical.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/parport.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/pktcdvd.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/pktcdvdN.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/port.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/power-switch.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/printer.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/ptmx.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/qftN.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/random.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/raw.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/rawctl.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/rfkill.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/rtc.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/sch.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/sclp_line.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/sgi.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/sonypi.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/sound-alsa-control.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/sound-audio.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/tty.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/tty3270.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/ttyN.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/ttyS.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/ttysclp.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/usb-bus.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/v4l.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/vcs.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/virtio.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/virtualbox-usb.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/virtualbox-usb_device.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/z90crypt.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/zzz-event-push.act
-rw-r--r--	t0b1nux/users	etc/vdev/actions/zzz-udev-compat.act
-rw-r--r--	t0b1nux/users	etc/vdev/ifnames.conf
-rw-r--r--	t0b1nux/users	etc/vdev/vdevd.conf
drwxr-xr-x	t0b1nux/users	lib/
-rwxr-xr-x	t0b1nux/users	lib/libudev.so.1
lrwxrwxrwx	t0b1nux/users	lib/libudev.so.1.5.2 -> libudev.so.1
drwxr-xr-x	t0b1nux/users	lib/vdev/
-rwxr-xr-x	t0b1nux/users	lib/vdev/VBoxCreateUSBNode.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/block.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/char.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/daemonlet
-rwxr-xr-x	t0b1nux/users	lib/vdev/dev-setup.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/disk.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/dm-disk.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/echo_n
-rwxr-xr-x	t0b1nux/users	lib/vdev/event-push.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/event-put
-rwxr-xr-x	t0b1nux/users	lib/vdev/firmware.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/gen-ifnames.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/hwdb-props.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/hwdb.sh
drwxr-xr-x	t0b1nux/users	lib/vdev/hwdb/
-rw-r--r--	t0b1nux/users	lib/vdev/hwdb/hwdb.squashfs
-rwxr-xr-x	t0b1nux/users	lib/vdev/ifname.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/input.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/modprobe.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/netdev.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/optical.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/permissions.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/power-switch.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/sound.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/stat_ata
-rwxr-xr-x	t0b1nux/users	lib/vdev/stat_bus
-rwxr-xr-x	t0b1nux/users	lib/vdev/stat_input
-rwxr-xr-x	t0b1nux/users	lib/vdev/stat_net
-rwxr-xr-x	t0b1nux/users	lib/vdev/stat_optical
-rwxr-xr-x	t0b1nux/users	lib/vdev/stat_path
-rwxr-xr-x	t0b1nux/users	lib/vdev/stat_scsi
-rwxr-xr-x	t0b1nux/users	lib/vdev/stat_usb
-rwxr-xr-x	t0b1nux/users	lib/vdev/stat_v4l
-rw-r--r--	t0b1nux/users	lib/vdev/subr-event.sh
-rw-r--r--	t0b1nux/users	lib/vdev/subr-hwdb.sh
-rw-r--r--	t0b1nux/users	lib/vdev/subr.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/sysfs.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/udev-compat.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/usb-name.sh
-rwxr-xr-x	t0b1nux/users	lib/vdev/v4l.sh
drwxr-xr-x	t0b1nux/users	sbin/
-rwxr-xr-x	t0b1nux/users	sbin/vdevd
-rwxr-xr-x	t0b1nux/users	sbin/vdevfs
drwxr-xr-x	t0b1nux/users	usr/
drwxr-xr-x	t0b1nux/users	usr/include/
-rw-r--r--	t0b1nux/users	usr/include/libudev.h
drwxr-xr-x	t0b1nux/users	usr/share/
drwxr-xr-x	t0b1nux/users	usr/share/pkgconfig/
-rw-r--r--	t0b1nux/users	usr/share/pkgconfig/libudev.pc
